# etc/systemd/user/rewards-max.timer
[Unit]
Description=Semi-day Rewards Max

[Timer]
OnBootSec=15min
OnUnitActiveSec=8h
