# etc/systemd/system/tmpfiles-cleanup.timer
[Unit]
Description=Daily Cleanup of Temporary Files

[Timer]
OnBootSec=15min
OnUnitActiveSec=1d
